# rootfs/etc/services.d/mysql/run
#!/usr/bin/with-contenv bash

/usr/bin/mysqld_safe